# cookbooks/enterprise_linux/templates/default/etc/profile.d/history.sh.erb
<%= node['linux']['chef']['file_header'] %>

umask <%= node['linux']['user']['umask'] %>

HISTTIMEFORMAT='%F %T '
readonly HISTTIMEFORMAT
export HISTTIMEFORMAT
HISTSIZE=5000
readonly HISTSIZE
export HISTSIZE
