# cookbooks/enterprise_linux/templates/default/etc/profile.d/prompt.sh.erb
<%= node['linux']['chef']['file_header'] %>

umask <%= node['linux']['user']['umask'] %>

if [ "$PS1" ]; then
  if [[ ${EUID} == 0 ]] ; then
    PS1='\[\033[01;31m\]\h:\[\033[01;34m\]\W \u\$\[\033[00m\] '
  else
    PS1='\[\033[01;32m\]\h:\[\033[01;34m\]\W \u\$\[\033[00m\] '
  fi
fi
