# cookbooks/enterprise_linux/templates/default/etc/profile.d/tmout.sh.erb
<%= node['linux']['chef']['file_header'] %>

umask <%= node['linux']['user']['umask'] %>

if [ ! "${TMOUT}" ]
then
  readonly TMOUT=<%= node['linux']['security']['shell_timeout'] %>
  export TMOUT
fi
